# staging/system-update.service
[Unit]
Description=Comprehensive System Update
Documentation=file:///home/chris/projects/human-oriented-environment/staging/comprehensive-update.py
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Set nice level for background updates
Nice=19

# Environment
Environment="HOME=/home/chris"
Environment="PATH=/home/chris/.local/bin:/home/chris/.cargo/bin:/home/chris/.nvm/versions/node/v24.6.0/bin:/usr/local/bin:/usr/bin:/bin"

# For SDKMAN
Environment="SDKMAN_DIR=/home/chris/.sdkman"

# Script location
ExecStart=/home/chris/projects/human-oriented-environment/staging/comprehensive-update.py

# Don't run if on battery (for laptops)
ExecStartPre=/bin/bash -c 'test $(cat /sys/class/power_supply/AC/online 2>/dev/null || echo 1) -eq 1'

# Log output
StandardOutput=append:/home/chris/.local/share/update-logs/systemd-update.log
StandardError=append:/home/chris/.local/share/update-logs/systemd-update-error.log

# Restart policy - don't restart on failure
Restart=no

# Timeout for the entire update process (2 hours should be plenty)
TimeoutStartSec=7200

# Kill only the main process, not child processes
KillMode=process

# Send SIGTERM first, then SIGKILL after 30s
TimeoutStopSec=30

[Install]
WantedBy=default.target