# staging/system-update.timer
[Unit]
Description=Run Comprehensive System Update Weekly
Documentation=man:systemd.timer(5)
Requires=system-update.service

[Timer]
# Run every Sunday at 3 AM
OnCalendar=weekly
OnCalendar=Sun *-*-* 03:00:00

# Also run 15 minutes after boot if we missed a scheduled run
OnBootSec=15min

# Randomize start time by up to 30 minutes to avoid load spikes
RandomizedDelaySec=30min

# If the system was powered down when the timer was scheduled, run it when powered on
Persistent=true

# Don't run multiple instances
Unit=system-update.service

[Install]
WantedBy=timers.target